# graphite/files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}
{%- set workers = grains.num_cpus * 2 + 1 %}
[program:{{ service_name }}]
command=/usr/local/bin/gunicorn_django --user graphite --group graphite --bind 127.0.0.1:8787 --log-file=/opt/graphite/storage/log/webapp/gunicorn.log --workers {{ workers }} --timeout 600 /opt/graphite/webapp/graphite/settings.py
#command=/usr/local/bin/gunicorn --user graphite --group graphite --bind 127.0.0.1:8787 --log-file=/opt/graphite/storage/log/webapp/gunicorn.log --workers {{ workers }} --timeout 600 /opt/graphite/conf/graphite.wsgi
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT
user=graphite
